# files/crontab/cronjob.alpine
1 0 * * * cd /var/www/kanboard && ./kanboard cronjob
