# files/services.d/worker/run
#!/bin/execlineb -P
s6-setuidgid nginx /var/www/app/kanboard worker